# Daemon/openrc/ydotoold.service-openrc.in
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Starts ydotoold service"
supervisor="supervise-daemon"
command="@YDOTOOLD_SERVICE_COMMAND@"

reload() {
	ebegin "Reloading $RC_SVCNAME"
	/usr/bin/kill -HUP $MAINPID
	eend $?
}
